# templates/gunicorn.service.erb
[Unit]
Description=<%= @name %> gunicorn app
After=network.target

[Service]
PIDFile=/run/gunicorn/<%= @name %>
User=<% @owner %>
Group=<%= @group %>
RuntimeDirectory=<%= @dir %>
WorkingDirectory=<%= @dir %>
ExecStart=/usr/bin/gunicorn -c /etc/gunicorn.d/<%= @name %> --pid /var/run/gunicorn/<%= @name %> 
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
PrivateTmp=true

[Install]
WantedBy=multi-user.target
